# puppet/modules/pagetty/templates/upstart/pagetty_crawler.conf.erb
description "pagetty crawler"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

setuid pagetty

script
  export NODE_ENV=production
  cd /srv/pagetty/current
  node /srv/pagetty/current/crawler.js 2>&1 | logger -t "crawler.js"
end script
